# systemd/.config/systemd/user/urxvtd.service
[Unit]
Description=urxvtd

[Service]
Type=forking
ExecStart=/usr/bin/urxvtd -q -o -f

[Install]
WantedBy=default.target
